# docker/s6-rc.d/mcp-server/run
#!/command/with-contenv sh

# MCP server provides Meticulous machine tools for external MCP clients (Claude Desktop, etc.)
# External clients may use their own LLM/API key — MeticAI's GEMINI_API_KEY is not required.
# Core profile creation uses the Python Gemini SDK directly (#214).

cd /app/mcp-server
export PYTHONPATH=/app/mcp-server/meticulous-mcp/src
export FASTMCP_HOST=0.0.0.0
export FASTMCP_PORT=8080
export METICULOUS_API_URL="http://${METICULOUS_IP:-meticulous.local}"
exec python run_http.py
