# docker/s6-rc.d/meticulous-bridge/run
#!/command/with-contenv sh
cd /app/bridge
exec python start_bridge.py
