# docker/s6-rc.d/server/run
#!/command/with-contenv sh
cd /app/server
exec uvicorn main:app --host 0.0.0.0 --port 8000 --workers 1
