# docker/s6-rc.d/mosquitto/run
#!/command/with-contenv sh
exec mosquitto -c /etc/mosquitto/mosquitto.conf
